# files/etc/service/nginx/run
#!/bin/sh
set -e

mkdir -p /var/run/nginx
chmod 0755 /var/run/nginx
nginx
