# files/etc/cron.daily/geoip2-update
#!/bin/bash
# download geolite2 file
cd /etc/nginx
./geoip2-download.sh
