# files/etc/cron.daily/geoip2-update
#!/bin/bash
# download geolite2 file
bash /etc/nginx/geoip2-download.sh

# reload nginx
nginx -s reload